# provisioning/base/files/runit/app/run
#!/bin/sh
echo "=============== STARTING APP SERVER ==================="
exec cd /app && npm run start
