# provisioning/base/files/runit/app/run
#!/bin/bash
echo "=============== STARTING APP SERVER ==================="
exec npm run --prefix /app generate && bpn run --prefix /app start
