# old/provisioning/base/files/runit/php-fpm/run
#!/bin/sh
set -e
echo "=============== STARTING PHP-FPM ==================="
exec /usr/local/sbin/php-fpm --nodaemonize -R